# rc.local
#!/bin/sh
/bin/mount -o remount,exec /tmp

[ ! -f /opt/docker-bbb/setup-executed ] && /opt/docker-bbb/setup.sh
