# rc.local
#!/bin/sh
# [ ! -f /opt/docker-bbb/setup-executed ] && /opt/docker-bbb/setup.sh
